# deploy/airgap/templates/edossier-app.service
[Unit]
Description=E-Dossier Application
After=network.target
ConditionPathExists=__APP_ROOT__/current/server.js

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_ROOT__/current
EnvironmentFile=__APP_ROOT__/shared/app.env
Environment=HOSTNAME=127.0.0.1
Environment=PORT=__APP_PORT__
ExecStart=__NODE_BIN__ server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
